# backend/gunicorn.service
[Unit]
Description=Gunicorn daemon for TrainingProtec Backend
After=network.target

[Service]
Type=exec
User=root
Group=root
WorkingDirectory=/var/www/tranning/backend
Environment="PATH=/var/www/tranning/backend/venv/bin"
EnvironmentFile=-/var/www/tranning/backend/.env

# Create required directories before starting
ExecStartPre=/bin/mkdir -p /var/log/gunicorn
ExecStartPre=/bin/mkdir -p /var/www/tranning/backend/static/uploads
ExecStartPre=/bin/mkdir -p /var/www/tranning/backend/instance

ExecStart=/var/www/tranning/backend/venv/bin/gunicorn --config /var/www/tranning/backend/gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
