# backend/gunicorn.service
[Unit]
Description=Gunicorn daemon for TrainingProtec Backend
After=network.target

[Service]
Type=exec
User=www-data
Group=www-data
WorkingDirectory=/var/www/tranning/backend
Environment="PATH=/var/www/tranning/backend/venv/bin"
EnvironmentFile=-/var/www/tranning/backend/.env
ExecStart=/var/www/tranning/backend/venv/bin/gunicorn --config /var/www/tranning/backend/gunicorn.conf.py app:app
ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=30
PrivateTmp=true
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
